# projects/workbench/com.ainary.daily-digest.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ainary.daily-digest</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/Users/florianziesche/.openclaw/workspace/projects/workbench/daily-digest-cron.sh</string>
    </array>
    
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>21</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    
    <key>StandardOutPath</key>
    <string>/Users/florianziesche/.openclaw/workspace/logs/daily-digest-stdout.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/florianziesche/.openclaw/workspace/logs/daily-digest-stderr.log</string>
    
    <key>RunAtLoad</key>
    <false/>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
